# old.frontend/AutomacaoCantina-master/service/controle_cantina.service
[Unit]
Description=Controle Cantina

[Service]
ExecStart=/home/pi/AutomacaoCantina/almoco/controleCantina.py
Restart=always

[Install]
WantedBy=default.target
